# infrastructure/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=appuser
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:uvicorn]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
directory=/app
autostart=true
autorestart=true
startsecs=10
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1

[program:groupchat-summary-worker]
command=python -m app.groupchat.summary.worker --loop --interval-seconds=%(ENV_GROUPCHAT_SUMMARY_POLL_SECONDS)s
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
killasgroup=true
priority=900
environment=PYTHONUNBUFFERED=1,GROUPCHAT_SUMMARY_POLL_SECONDS=5,GROUPCHAT_SUMMARY_LOG_LEVEL=INFO

[program:groupchat-followup-worker]
command=python -m app.groupchat.followup.worker --loop --interval-seconds=%(ENV_GROUPCHAT_FOLLOWUP_POLL_SECONDS)s
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
killasgroup=true
priority=910
environment=PYTHONUNBUFFERED=1,GROUPCHAT_FOLLOWUP_POLL_SECONDS=10,GROUPCHAT_FOLLOWUP_LOG_LEVEL=INFO

[program:daily-email-worker]
command=python -m app.proactive.email.worker --loop --interval-seconds=%(ENV_DAILY_EMAIL_POLL_SECONDS)s
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
killasgroup=true
priority=930
environment=PYTHONUNBUFFERED=1,DAILY_EMAIL_POLL_SECONDS=300,DAILY_EMAIL_LOG_LEVEL=INFO

[program:proactive-outreach-worker]
command=python -m app.proactive.outreach.worker --loop --interval-seconds=%(ENV_PROACTIVE_OUTREACH_POLL_SECONDS)s
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopwaitsecs=30
killasgroup=true
priority=940
environment=PYTHONUNBUFFERED=1,PROACTIVE_OUTREACH_POLL_SECONDS=300,PROACTIVE_OUTREACH_LOG_LEVEL=INFO
